# warpgui/all/etc/systemd/user/warp-update.service
[Unit]
Description=WARP Update Service

[Service]
Type=simple
ExecStart=warp-update
KillMode=mixed

[Install]
WantedBy=default.target
